# ocmauthd.service
[Unit]
Description=Authentication daemon for CERNBox
After=syslog.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/log/ocmauthd
ExecStart=/usr/local/bin/ocmauthd
StandardOutput=null
StandardError=syslog

[Install]
WantedBy=multi-user.target

